# debian/prerm
#!/bin/sh
set -e
PATH=/bin:/usr/bin:/usr/sbin

#DEBHELPER#

rm -rf /var/lib/autolog

invoke-rc.d autolog stop

exit 0
